# roles/kafka/templates/supervisor.conf
[program:kafka]
environment=JAVA_HOME={{java_home}}
command={{install_dir}}/kafka/bin/kafka-server-start.sh {{install_dir}}/kafka/config/server.properties
user={{ user }}

