# roles/hbase-master/templates/supervisor.conf
[program:hbase-master]
environment=JAVA_HOME={{java_home}},HADOOP_CONF_DIR={{install_dir}}/hadoop/etc/datanode
command={{install_dir}}/hbase/bin/hbase master start
user={{ user }}

