# roles/nodemanager/templates/supervisor.conf
[program:nodemanager]
environment=JAVA_HOME={{java_home}}
command={{install_dir}}/hadoop/bin/yarn --config {{install_dir}}/hadoop/etc/nodemanager nodemanager
user=root

