# roles/zeppelin2/templates/supervisor.conf
[program:zeppelin2]
environment=JAVA_HOME={{java_home}},ZEPPELIN_PORT=8787
command={{ install_dir }}/zeppelin2/bin/zeppelin.sh
user=root
