# roles/hbase-regionserver/templates/supervisor.conf
[program:hbase-regionserver]
environment=JAVA_HOME={{java_home}},HADOOP_CONF_DIR={{install_dir}}/hadoop/etc/datanode
command={{install_dir}}/hbase/bin/hbase regionserver start
user={{ user }}

