# roles/zookeeper/templates/supervisor.conf
[program:zookeeper]
environment=JAVA_HOME={{java_home}}
command={{install_dir}}/zookeeper/bin/zkServer.sh start-foreground
user=root

