# roles/hdfs-datanode-ha/templates/supervisor.conf
[program:datanode]
environment=JAVA_HOME={{java_home}}
command={{install_dir}}/hadoop/bin/hdfs --config {{install_dir}}/hadoop/etc/datanode datanode 
user=root

