# roles/hbase-rest/templates/supervisor.conf
[program:hbase-rest]
environment=JAVA_HOME={{java_home}}
command={{install_dir}}/hbase/bin/hbase rest start -p 5656
user={{ user }}

